# config/pi-gen/stage3/01-microscopi/files/microscopi.service
[Unit]
Description=Microscopi Camera Application
After=local-fs.target
Wants=local-fs.target

[Service]
User=microscopi
Group=microscopi
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_AUDIODRIVER=dummy
Environment=HOME=/home/microscopi
ExecStart=/usr/local/bin/microscopi
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
