# config/pi-gen/stage3/01-microscopi/files/microscopi-wifi-ap.service
[Unit]
Description=Microscopi WiFi AP fallback
Documentation=https://github.com/microscopi/microscopi
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/microscopi-wifi-ap-check.sh

[Install]
WantedBy=multi-user.target
